# debian/asl-telemetry.timer
[Unit]
Description=Run asl-telemetry

[Timer]
OnBootSec=120
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
