# share/asl-play-arn.timer
[Unit]
Description=PLay ARN on the schedule
Requires=asl-play-arn.service

[Timer]
OnCalendar=Tue *-*-* 19:30:00

[Install]
WantedBy=timers.target

